# systemd/fw_amend.timer
[Unit]
Description=Firewall check and updates periodically

[Timer]
OnCalendar = *-*-* *:12,32,52

[Install]
WantedBy=timers.target
